# collection/roles/postgresql-proxy/files/postgresql-proxy.socket
[Unit]
Description=Socket activation for Postgresql service

[Socket]
ListenStream=15432

[Install]
WantedBy=sockets.target
